# snmptrapd_custom/supervisor/conf.d/snmptrapd.conf
[program:snmptrapd]
#directory=/opt/bgpalert
command=/usr/sbin/snmptrapd -Lo -f
autostart=true
autorestart=true
stderr_logfile=/var/log/snmptrapd.err.log
stdout_logfile=/var/log/snmptrapd.out.log
